# docker/complement/conf/postgres.supervisord.conf
[program:postgres]
# Disable crash-safety features for speed: fsync, synchronous commit, full page writes
command=/usr/local/bin/prefix-log gosu postgres postgres
    -c fsync=off
    -c synchronous_commit=off
    -c full_page_writes=off
    -c wal_level=minimal
    -c max_wal_senders=0

# Only start if START_POSTGRES=true
autostart=%(ENV_START_POSTGRES)s

# Lower priority number = starts first
priority=1

autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Use 'Fast Shutdown' mode which aborts current transactions and closes connections quickly.
# (Default (TERM) is 'Smart Shutdown' which stops accepting new connections but
#  lets existing connections close gracefully.)
stopsignal=INT
